# deploy/systemd/quietgrid-web.service
[Unit]
Description=QuietGrid Web Monitor
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/quietgrid
ExecStart=/home/ubuntu/quietgrid/.venv/bin/python web.py
Restart=always
RestartSec=10
EnvironmentFile=/home/ubuntu/quietgrid/.env

[Install]
WantedBy=multi-user.target
